# webserver/templates/default/mail-on-restart.conf.erb
description "Job that sends an e-mail when a service is restarted"
author  "Mark Desnoyer <desnoyer@neon-lab.com>"

start on stopped <%= @service %> RESULT=failed

script

echo "<%= @service %> crashed on <%= @host %> at `date`. It has probably restarted. The last logs are:\n `tail -n 100 <%= @log_file %>`" | mail -s "<%= @service %> crashed on <%= @host %>" -aFrom:<%= @host%>@neon-lab.com -areply-to:noreply@neon-lab.com <%= @email %>


end script