# webserver/templates/default/system_metrics_service.conf.erb
description "Neon system metrics"

start on runlevel [234]
stop on runlevel [0156]

kill timeout 30
limit nofile 65536 65536

setuid <%= @user %>
setgid <%= @group %>

chdir <%= @neon_root_dir %>
exec <%= @neon_root_dir %>/monitoring/check_system_health.py --carbon_server <%= @carbon_host %> --carbon_port <%= @carbon_port %>
respawn
